# server/deploy/systemd/codex-runtime-health-check.timer
[Unit]
Description=Daily host-side Codex runtime health check

[Timer]
OnCalendar=*-*-* 05:20:00
RandomizedDelaySec=15min
Persistent=true
Unit=codex-runtime-health-check.service

[Install]
WantedBy=timers.target
